# crontab
* * * * * root php /var/www/index.php verbose=1
